# data/service/hyelicht.service.in
# SPDX-License-Identifier: GPL-2.0-only OR GPL-3.0-only OR LicenseRef-KDE-Accepted-GPL
# SPDX-FileCopyrightText: 2021-2022 Eike Hein <sho@eikehein.com>

[Unit]
Description=hyelicht Controller

[Service]
Type=simple
ExecStart=@KDE_INSTALL_FULL_BINDIR@/hyelicht --onboard
Environment=QT_QPA_PLATFORM=eglfs
Environment=QT_QPA_EGLFS_INTEGRATION=eglfs_kms
Environment=QT_QPA_EGLFS_WIDTH=1024
Environment=QT_QPA_EGLFS_HEIGHT=600
Environment=QT_QPA_EGLFS_ALWAYS_SET_MODE=1
Environment=QT_QPA_EGLFS_KMS_ATOMIC=1
Environment=QT_LOGGING_RULES=com.hyerimandeike.*=true

[Install]
WantedBy=default.target
